# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:serf]
command=/usr/local/bin/serf agent -config-dir=/etc/serf
