# back-end/bundle/supervisord.conf
[supervisord]
nodaemon=true

[program:serf]
command=/usr/local/bin/serf agent -config-dir=/etc/serf
